# dmz/490dmz.service
[Unit]
Description=SystemD Rabbit Listener

[Service]
Type=simple
ExecStart=/usr/bin/php -f /home/hsm24mq/git/dmz/apiServer.php

[Install]
WantedBy=multi-user.target
